# slash/etc/profile.d/spark.sh
export SPARK_HOME=/opt/spark
export SPARK_SERVER=
export SPARK_HOST=
export SPARK_MASTER=
export PATH=/opt/spark/bin:${PATH}
