# slash/opt/spark/extras/spark.service
[Unit]
Description=Spark Service

[Service]
ExecStart=/bin/bash -c "/opt/spark/extras/spark-start"

[Install]
WantedBy=multiuser.target
